# configs/systemd/lobby-kiosk.target
[Unit]
Description=Lobby Kiosk System
Wants=lobby-app.service lobby-display.service lobby-watchdog.service

[Install]
WantedBy=multi-user.target